# dist/kahlu-web.service
[Unit]
Description=kahlu-web
After=network.target

[Service]
Type=simple
User=kahlu
WorkingDirectory=/home/kahlu/live
Environment="RAILS_ENV=production"
Environment="PORT=3000"
ExecStart=/home/kahlu/.rbenv/shims/bundle exec puma -C config/puma.rb
ExecReload=/bin/kill -SIGUSR1 $MAINPID
TimeoutSec=15
Restart=always

[Install]
WantedBy=multi-user.target
